# package/deb/DEBIAN/prerm
#!/bin/bash
set -e

# Update services
deb-systemd-invoke stop backlight-tracer || systemctl stop backlight-tracer
deb-systemd-helper disable backlight-tracer

# Done
exit 0
